# tt-rss/service/update/run
#!/bin/sh

exec chpst -u www-data "$TTRSS_PHP_EXECUTABLE" /usr/local/tt-rss/html/update_daemon2.php
