# tt-rss/service/update/run
#!/bin/sh

exec chpst -u www-data /usr/local/bin/php /usr/local/tt-rss/update_daemon2.php