# repman/service/messenger/run
#!/bin/sh

exec chpst -u www-data /usr/local/bin/repman messenger:consume async --limit 50
