# backuppc/service/backuppc/run
#!/bin/sh

exec chpst -u "$BACKUPPC_USER" /usr/local/BackupPC/bin/BackupPC
