# backuppc/service/fcgiwrap/run
#!/bin/sh

rm -rf "$BACKUPPC_FCGI"

exec spawn-fcgi -n -s "$BACKUPPC_FCGI" -u "$BACKUPPC_UUID" -g 82 -M 0770 -- /usr/bin/fcgiwrap
